# orange-pi-daemon/sms-dashboard-daemon.service
[Unit]
Description=SMS Dashboard Daemon
After=network-online.target ModemManager.service
Wants=network-online.target
Requires=ModemManager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sms-dashboard-daemon
Restart=always
RestartSec=10
User=sms-daemon
Group=dialout

# Environment
Environment="SMS_API_URL=https://sexy.qzz.io"
# Set API key via override or EnvironmentFile
# Environment="SMS_API_KEY=your-api-key"
EnvironmentFile=-/etc/sms-dashboard-daemon/config

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sms-dashboard

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sms-dashboard-daemon

[Install]
WantedBy=multi-user.target